# satori.cribfinder/deboot_builder/kernel/template_oneiric/initrd/etc/udev/rules.d/80-compcache.rules
KERNEL=="zram0", ACTION=="add", \
  RUN+="/sbin/compcache-enable zram 50% zram0"
