# satori.judge/runner/debian/runner.upstart
description	"Runner daemon"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec runnerd
