# satori-testing/debian/satori-testing.upstart
description	"Satori run daemon"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec satori_rund
